# supervisor/simulator_service.conf
[program:simulator_service]
directory=/root/ai-voice-agent/services/simulator_service
command=/root/ai-voice-agent/services/simulator_service/simulator_venv/bin/python simulate_complete_system.py
autostart=false
autorestart=false
startsecs=0
exitcodes=0
stdout_logfile=/var/log/ai_agent/simulator_service.out.log
stderr_logfile=/var/log/ai_agent/simulator_service.err.log
